# deploy/systemd/sync-api.service
[Unit]
Description=Sync API
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
User=sync
Group=sync
WorkingDirectory=/opt/sync/backend
EnvironmentFile=/etc/sync/backend.env
ExecStart=/opt/sync/backend/.venv/bin/uvicorn sync_backend.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
